# conf/mariadb.conf
[supervisord]
nodaemon=true

[program:hcs]
command=/bin/bash /usr/bin/hcs-start

[program:mysql]
command=service mysql start
